# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:node_vless]
command=node /app/vless_server.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:node_keepalive]
command=node /app/keepalive.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:cloudflared]
command=/app/cloudflared tunnel --no-autoupdate run --token eyJhIjoiYWRkMWM4MzNjZWY4YjRlY2ZhZjUzNmJlNTlhYmNkYjgiLCJ0IjoiNWNiMmE3YjEtYThjYy00ODYxLTk5YjItOGNmZTdkN2U4YTFlIiwicyI6Ik5Ua3hORE0yWldZdE56VmhZaTAwTVRVeExXSmxOREl0WWpKaU9URmtZV1prT1RZMiJ9
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

# 如果你不需要 sing-box，这一段可以在配置里删掉以节省内存。
# 因为你手写的 Node 脚本已经替代了 sing-box 入口的功能。
[program:singbox]
command=/app/sing-box run -c /app/config.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
